# mirror/mirrordate.timer
[Unit]
Description=Update TIME file on abf-downloads

[Timer]
OnBootSec=5min
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
