# motd/motd.timer
[Unit]
Description=Run motd.service every 10 minutes

[Timer]
OnBootSec=5min
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
