# sing-rpms/openmandriva-sign-rpm.timer
[Unit]
Description=Sign OpenMandriva RPMS trigger

[Timer]
OnBootSec=5min
OnUnitActiveSec=3m

[Install]
WantedBy=timers.target
